# data/goa-daemon.service.in
[Unit]
Description=GNOME Online Accounts service
# Terminate automatically on logout from GUI session
PartOf=graphical-session.target

[Service]
ExecStart=@libexecdir@/goa-daemon
Type=dbus
BusName=org.gnome.OnlineAccounts
TimeoutStopSec=5
