# data/goa-identity-service.service.in
[Unit]
Description=GNOME identity service for Kerberos
# Terminate automatically on logout from GUI session
PartOf=graphical-session.target

[Service]
ExecStart=@libexecdir@/goa-identity-service
Type=dbus
BusName=org.gnome.Identity
TimeoutStopSec=5
